# api/core/v1beta1/_assets/etc/profile.d/99-ptd.sh
# NOTE: This file is injected by rstudio/ptd team-operator
umask 077
